# target/linux/ipq40xx/chromium/base-files/etc/hotplug.d/ieee80211/99-google-wifi-steering
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0

. /lib/functions.sh

[ "$(board_name)" = "google,wifi" ] || exit 0

/etc/init.d/packet_steering reload
